# agent-support/xcode/com.gitai.xcode-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gitai.xcode-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/PATH/TO/git-ai-xcode-watcher</string>
        <string>--path</string>
        <string>/PATH/TO/YOUR/PROJECT</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/git-ai-xcode-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/git-ai-xcode-watcher.log</string>
</dict>
</plist>
